# pi/wled-scheduler-update.path
[Unit]
Description=Watch for WLED Scheduler update flag

[Path]
PathExists=/opt/wled-scheduler/data/update.flag

[Install]
WantedBy=multi-user.target
